# incron/files/incron.init
#!/bin/sh /etc/rc.common
#
# Author: Zhifeng Gu <guzhifeng1979@hotmail.com>
# Date: 2012/05
# License: GPL

START=99
STOP=00

start() {
    incrond
}

stop() {
    incrond -k
}
